# nibe-api.service
[Unit]
Description=Nibe Autotuner API V4
After=network.target

[Service]
Type=simple
User=peccz
WorkingDirectory=/home/peccz/AI/nibe_autotuner
Environment="PATH=/home/peccz/AI/nibe_autotuner/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/peccz/AI/nibe_autotuner/src"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/peccz/AI/nibe_autotuner/venv/bin/python src/api_server.py

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
